# providerscripts/webserver/configuration/application/lighttpd/rc.local
#!/bin/bash
/usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf

/bin/sleep 10

lighties="`/usr/bin/ps -ef | /bin/grep lighttp | /bin/grep -v grep | /usr/bin/wc -l`"

while ( [ "${lighties}" = "0" ] )
do
    /usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf
    lighties="`/usr/bin/ps -ef | /bin/grep lighttp | /bin/grep -v grep | /usr/bin/wc -l`"
    /bin/sleep 10
done

exit 0
